# systemd/vigiclient.service.d/encode.conf
[Service]
Environment=VIGI_USE_FFMPEG=1
Environment=LD_LIBRARY_PATH=/usr/hobot/lib:/usr/lib
Environment=VIGI_DUMP_PATH=/tmp/vigi_live.h264
